# clamir-network-ip-apply.path
[Path]
PathModified=/mnt/mmc/sys/secondIP

[Install]
WantedBy=default.target