# clamir-scc-util-convert-table-monitor.path
[Path]
PathModified=/mnt/mmc/sys/table.nuc

[Install]
WantedBy=default.target
